# crontab
0 12 * * * /bin/systemctl restart emojicode-httpd.service